# etc/profile.d/alfresco_aliases.sh

alias AGP='vim /opt/alfresco/tomcat/shared/classes/alfresco-global.properties'

alias COUT='less /opt/alfresco/tomcat/logs/catalina.out'

alias MMT='/opt/alfresco/java/bin/java -jar /opt/alfresco/bin/alfresco-mmt.jar'
